# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch
autorestart=true

[program:cron]
command=/usr/sbin/cron -f
autorestart=true
